# deploy/graphenmail.service
[Unit]
Description=GraphenMail Email Extractor
After=network.target

[Service]
Type=simple
User=graphenmail
Group=graphenmail
WorkingDirectory=/opt/graphenmail
EnvironmentFile=-/etc/graphenmail/env
StateDirectory=graphenmail
ConfigurationDirectory=graphenmail
ExecStart=/opt/graphenmail/venv/bin/gunicorn --bind 127.0.0.1:5000 --workers 1 --threads 8 --timeout 300 "web:create_app()"
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment
Environment=PYTHONUNBUFFERED=1

# Hardening
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
NoNewPrivileges=yes
ReadWritePaths=/var/lib/graphenmail /etc/graphenmail /opt/graphenmail

[Install]
WantedBy=multi-user.target
